# deploy/newt.service
[Unit]
Description=Newt Connector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/algodomoiot/newt.env
ExecCondition=/bin/sh -c '[ "${NEWT_ENABLED:-0}" = "1" ] && [ -n "${NEWT_ID:-}" ] && [ -n "${NEWT_SECRET:-}" ] && [ -n "${NEWT_ENDPOINT:-}" ]'
ExecStart=/usr/local/bin/newt --id ${NEWT_ID} --secret ${NEWT_SECRET} --endpoint ${NEWT_ENDPOINT}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
